# backend/aud_gateway/flexepin.socket
[Unit]
Description=Flexepin Socket

[Socket]
ListenStream=/run/flexepin.socket
Accept=yes

[Install]
WantedBy=sockets.target